# onedrive-lumen.service
[Unit]
Description=OneDrive sync (download-only) for Lumen
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User=admin
WorkingDirectory=/home/admin
Environment=HOME=/home/admin
Environment=XDG_CONFIG_HOME=/home/admin/.config
ExecStart=/usr/local/bin/onedrive --confdir /home/admin/.config/onedrive --sync --download-only --verbose
Nice=10
SuccessExitStatus=0

[Install]
WantedBy=multi-user.target
